# debian/hwmonitor.cron.daily
#!/bin/sh

# restart the hwmonitor service
stop hwmonitor
start hwmonitor
